# xmonad/profile.d/ghcup.sh
export PATH=$PATH:/opt/ghcup/.ghcup/bin
export PATH=$PATH:/opt/ghcup/.cabal/bin
export GHC_PACKAGE_PATH=/opt/ghcup/.cabal/store/ghc-8.10.7/package.db:
